# ollama.service
[Unit]
Description=Ollama LLM Server with Dual GPU Support
After=network.target

[Service]
Type=simple
User=ollama
Environment="OLLAMA_MODELS=/var/lib/ollama/models"
Environment="CUDA_VISIBLE_DEVICES=0,1"
ExecStart=/usr/bin/ollama serve
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

[Install]
WantedBy=multi-user.target
