# raptorcare.service
[Unit]
Description=RaptorCare FastAPI Server
After=network.target postgresql.service

[Service]
Type=notify
User=raptorcare
WorkingDirectory=/opt/raptorcare
Environment="PATH=/opt/raptorcare/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/raptorcare/venv/bin/uvicorn server.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
